# BTTracker/udpt.service
[Unit]
Description=Simple UDP BitTorrent tracker written in C#.

[Service]
Type=notify
ExecStart=/usr/sbin/testapp

[Install]
WantedBy=multi-user.target